# crontab
*/2 * * * * cd /home/<user>/ && ./getinfo.sh
